# fabric_templates/gunicorn.conf
[program:%(project)s]
command=%(start_gunicorn_path)s
directory=%(project_dir)s
user=%(user)s
autostart=true
autorestart=true
redirect_stderr=True
